# scripts/setup_systemd/netsetup.service
[Unit]
Description = iproute2 auto start daemon

[Service]
ExecStart = /opt/netsetup.sh
Type=oneshot

[Install]
WantedBy = multi-user.target
